# recipes-core/udev-rules-vdo-terminal/udev-rules-vdo-terminal/99-vdo.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="1eaf", ATTRS{idProduct}=="0107", RUN+="@BINDIR@/udev-cp210x-newid.sh 1eaf 0107"
